# sshd/app/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
startsecs=3
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
